# nginx/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700
username = dummy
password = dummy

[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
username = dummy
password = dummy

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
autostart=true
startretries=0

[program:uwsgi]
command=uwsgi --ini /code/uwsgi.ini
stdout_logfile=/var/log/cache-db/uwsgi.log
stderr_logfile=/var/log/cache-db/uwsgi.log

[program:celery]
command=celery -A get_groups.groups_celery worker --loglevel=INFO
stdout_logfile=/var/log/cache-db/worker.log
stderr_logfile=/var/log/cache-db/worker.log

[program:celery-beat]
command=celery -A get_groups.groups_celery beat --loglevel=INFO -s /tmp/celerybeat-schedule --pidfile="/tmp/celerybeat.pid"
numprocs=1
stdout_logfile=/var/log/cache-db/beat.log
stderr_logfile=/var/log/cache-db/beat.log
autostart=false
autorestart=false
startsecs=10
